# supervisor.ini
[program:sshd]
command=/usr/sbin/sshd -f /etc/ssh/sshd_config